# rc
#
# PROVIDE: lambdabot
# REQUIRE: NETWORKING
# KEYWORD: shutdown
#

. /etc/rc.subr

name="lambdabot"
rcvar="lambdabot_enable"
desc="Bot for the Lambda Utility mod"

load_rc_config $name

: ${lambdabot_enable:=NO}
: ${lambdabot_user:=_lambdabot}
: ${lambdabot_group:=_lambdabot}
: ${lambdabot_command:=/usr/local/bin/lambdabot}
: ${lambdabot_flags:=}
: ${lambdabot_var:=/var/db/lambdabot}
: ${lambdabot_logdir:=/var/log/lambdabot}

pidfile="/var/run/lambdabot/lambdabot.pid"

command="${lambdabot_command}"

start_precmd="lambdabot_precmd"
start_cmd="lambdabot_start"

lambdabot_precmd()
{
    install -d -m 700 -o ${lambdabot_user} -g ${lambdabot_group} \
        ${lambdabot_var} ${lambdabot_logdir} /var/run/lambdabot/
}

lambdabot_start()
{
    /usr/bin/su -m "${lambdabot_user}" -c "${lambdabot_command} ${lambdabot_flags} >> ${lambdabot_logdir}/lambdabot.log 2>&1 &"

    pid=$(pgrep -f "${lambdabot_command}")
    if [ -n "$pid" ]; then
        echo $pid > ${pidfile}
        echo "pid: $pid"
    else
        echo "Error: lambdabot failed to start"
        return 1
    fi
}

lambdabot_stop() {
    if [ -f ${pidfile} ]; then
        pid=$(cat ${pidfile})
    if kill -TERM $pid; then
      echo "lambdabot stopped successfully"
      rm -f ${pidfile}
    else
      echo "unable to stop lambdabot"
    fi
    else
        echo "pid file not found, is lambdabot running?"
    fi
}


run_rc_command "$1"